# osocam/data/tmp/DzOSat/demonisat/lib/systemd/system/emu-OSCam_11718-r798.service
[Unit]

[Service]
Type=forking
ExecStart=/usr/bin/OSCam_11718-r798 &
NonBlocking=true
NotifyAccess=main

[Install]
WantedBy=multi-user.target
