# osocam/data/tmp/DzOSat/dreamelite2/lib/systemd/system/emu-OSCam_11718-r798.socket
[Unit]

[Socket]
ListenStream=/var/run/OSCam_11718-r798.pid
Backlog=8

[Install]
WantedBy=sockets.target
